# zabbix-alerts/profile.d/example.json.save
{
  "name": "example",
  "channels": [ 
    { "type": "zabbix2opsgenie",
      "enable": true
    },
    { "type": "opsgenie",
      "enable": false,
      "api": { 
        "key": "xxxxxxxx-xxxx-xxxx-xxxx-xxxxxxxxxxxx",
        "url": "https://api.eu.opsgenie.com/v1/json/zabbix"
      }
    }
  ]
}
